# image/arch-rpi4/rootfs-overlay/usr/lib/systemd/system/matterhub-maintenance.target
[Unit]
Description=MatterHub maintenance mode
Requires=multi-user.target
AllowIsolate=yes
