# image/arch-rpi4/rootfs-overlay/usr/lib/systemd/system/matterhub-bootstrapd.service
[Unit]
Description=MatterHub privileged bootstrap boundary
After=local-fs.target systemd-tmpfiles-setup.service iwd.service systemd-networkd.service
Wants=systemd-tmpfiles-setup.service iwd.service systemd-networkd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/lib/matterhub/bootstrapd.sh ensure-runtime

[Install]
WantedBy=multi-user.target
