# image/arch-rpi4/profiles/hosyond-5-dsi/systemd/matterhub-firstboot.service.d/10-display-profile.conf
[Unit]
Requires=matterhub-display-verify.service
After=matterhub-display-verify.service
