# image/arch-rpi4/profiles/hosyond-5-dsi/systemd/matterhub-maintenance.target
[Unit]
Description=MatterHub maintenance mode
Documentation=file:///usr/share/doc/matterhub/hosyond-5-dsi/README.md
AllowIsolate=yes
